# skeleton/etc/inittab
# Basic inittab.  As a busybox inittab this is completely different from the
# familiar sysvinit inittab.  Entries are:
#
#   <tty>:<ignored>:<action>:<process>
#
#   <tty>       Controlling tty for specified process, set to console if blank
#   <ignored>   This field is normally used for the runlevel, completely
#               ignored by busybox!
#   <action>    One of sysinit, respawn, askfirst, wait, once, restart,
#               ctrlaltdel, shutdown.
#   <process>   Command to run

::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty 9600 ttyS0
::once:/etc/init.d/rc

::shutdown:/etc/init.d/shutdown
::restart:/sbin/init
